# deploy/systemd/dyad.service
[Unit]
Description=Dyad Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/dyad
EnvironmentFile=/etc/dyad/dyad.env
ExecStart=/usr/bin/npm run web:server
Restart=on-failure
RestartSec=5
User=dyad
Group=dyad
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
